# cron/cron.d
0 5 * * * hpcstats hpcstats --batch-mode import --cluster=all
