# debian/cron.d
0 5 * * * hpcstats hpcstats import --cluster=all
