# debian/hpcstats-fsquota-agent.cron.d
0 0 * * * hpcstats /usr/share/hpcstats/bin/fsquota
