# debian/hpcstats-fsquota-agent.cron.d
0 0 * * * root /usr/share/hpcstats/bin/fsquota
