# debian/hpcstats-jobstats-launcher.cron.d
0 12 * * * hpcstats /usr/share/hpcstats/bin/launch-jobstats
